# vms/asterisk/lib/systemd/system/voicebot-dashboard.service
[Unit]
Description=Voicebot observability dashboard (read-only FastAPI + Tabler)
After=network.target

[Service]
Type=simple
User=asterisk
Group=asterisk
WorkingDirectory=/opt/voicebot-dashboard
ExecStart=/opt/voicebot-dashboard/venv/bin/uvicorn main:app --app-dir app --host ${VOICEBOT_DASHBOARD_BIND} --port ${VOICEBOT_DASHBOARD_PORT}
Restart=on-failure
RestartSec=5

Environment=VOICEBOT_DASHBOARD_BIND=127.0.0.1
Environment=VOICEBOT_DASHBOARD_PORT=8099
Environment=VOICEBOT_DASHBOARD_REFRESH_S=5
Environment=VOICEBOT_COMMON_DIR=/opt/voicebot-dashboard/common
EnvironmentFile=-/etc/asterisk-lab/env

# Filesystem isolation: the dashboard only ever reads voicebot state and the
# recordings spool; it never writes to either.
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadOnlyPaths=/var/lib/voicebot /var/spool/asterisk/monitor

# Kernel + control-plane isolation.
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true

# Process and capability restrictions.
NoNewPrivileges=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
